# recipes-multimedia/pulseaudio/pulseaudio-system/pulseaudio-system.service.in
[Unit]
Description=Sound Service

[Service]
Type=notify
ExecStart=@BINDIR@/pulseaudio --daemonize=no --exit-idle-time=-1 --disallow-exit=true -n -F @SYSCONFDIR@/pulse/default.pa --system
Restart=always

[Install]
WantedBy=multi-user.target
