# recipes-system/dr-runlevel/dr-runlevel/datarespons-apps.target.in
[Unit]
Description=datarespons applications target
Requires=@target@
Conflicts=rescue.service rescue.target
After=@target@ rescue.service rescue.target
AllowIsolate=yes
